# deploy/supervisor.conf
; /etc/supervisor/conf.d/trademinds.conf

[program:trademinds-backend]
command=/home/trademinds/trademinds/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2 --loop asyncio
directory=/home/trademinds/trademinds/backend
user=trademinds
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/trademinds/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/trademinds/backend_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/home/trademinds",PATH="/home/trademinds/trademinds/backend/venv/bin:%(ENV_PATH)s"
priority=100

[group:trademinds]
programs=trademinds-backend
